# root/etc/s6-overlay/s6-rc.d/sabnzbd/run
#!/command/with-contenv sh

. /opt/sabenv/bin/activate

exec s6-setuidgid appuser python3 /opt/sabnzbd/SABnzbd.py -s 0.0.0.0:8080 --config-file=/config/sabnzbd.ini
